# infra/deploy/mediamtx.service
[Unit]
Description=MediaMTX (conclave media server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/clawd-conclave
# EnvironmentFile provides YOUTUBE_STREAM_KEY + MEDIAMTX_PUBLISH_USER/PASS.
EnvironmentFile=/home/ubuntu/clawd-conclave/.env.stream
# Render the runtime config (substitutes publisher creds) before starting.
# MediaMTX v1.17 doesn't reliably do env substitution in the config itself.
ExecStartPre=/usr/bin/env bash /home/ubuntu/clawd-conclave/infra/render-mediamtx.sh
ExecStart=/usr/local/bin/mediamtx /home/ubuntu/clawd-conclave/infra/mediamtx.runtime.yml
Restart=on-failure
RestartSec=3
# Tighter resource caps so a runaway MediaMTX can't starve zkllmapi.
CPUQuota=150%
MemoryMax=2G

[Install]
WantedBy=multi-user.target
